# indiallsky_watchdog.service
[Unit]
Description=Indiallsky Watchdog Service
After=network.target

[Service]
User=ubuntu
WorkingDirectory=/home/ubuntu/indiallsky_watchdog
ExecStart=/home/ubuntu/indiallsky_watchdog/watchdog.py
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=indiallsky-watchdog
Restart=always

[Install]
WantedBy=multi-user.target
